# logstash-kibana-elasticsearch/supervisord.d/elasticsearch.ini
[program:elasticsearch]
command=/bin/bash -c "/etc/init.d/elasticsearch start"
autostart=true
